# supervisor.conf
[program:kiwiirc]
user=root
command=/root/KiwiIRC/kiwi -f start
redirect_stderr=true
autorestart=true
